# nginx-request_id-demo/ruby-backend/etc/supervisor/conf.d/supervisord.conf
[supervisord]
nodaemon=true

[program:splunkforwarder]
command=/opt/splunkforwarder/bin/splunk start --accept-license --nodaemon --no-prompt --answer-yes

[program:ruby]
command=bin/rails server -p 4001 -b 0.0.0.0
