# nginx-request_id-demo/nginx-frontend/etc/supervisor/conf.d/supervisord.conf
[supervisord]
nodaemon=true

[program:splunkforwarder]
command=/opt/splunkforwarder/bin/splunk start --accept-license --nodaemon --no-prompt --answer-yes

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"

[program:amplify-agent]
command=/etc/init.d/amplify-agent start
