# nginx-request_id-demo/go-backend/etc/supervisor/conf.d/supervisord.conf
[supervisord]
nodaemon=true

[program:splunkforwarder]
command=/opt/splunkforwarder/bin/splunk start --accept-license --nodaemon --no-prompt --answer-yes

[program:go]
command=go run /usr/local/go/src/hello-world/main.go

[program:amplify-agent]
command=/etc/init.d/amplify-agent start
